# deploy/com.igor.hermes-litellm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.igor.hermes-litellm</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/igorganapolsky/workspace/git/igor/hermes-eval/litellm/start-proxy.sh</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/igorganapolsky/.hermes/litellm-logs/proxy.out.log</string>
  <key>StandardErrorPath</key><string>/Users/igorganapolsky/.hermes/litellm-logs/proxy.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/Users/igorganapolsky/.local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>LITELLM_PORT</key><string>4010</string>
  </dict>
</dict>
</plist>
